# debian/repcpd.service
[Unit]
Description=Port Control Protocol (PCP) Daemon
After=network-online.target

[Service]
ExecStart=/usr/sbin/repcpd -n -f /run/repcpd/repcpd.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
